# myweb_scraper.service.example
# Legacy systemd unit (replaced by Docker — see DOCKER.md).
# If you still use systemd, install with:
#   sudo cp myweb_scraper.service.example /etc/systemd/system/myweb_scraper.service
#
# Typical pattern: run once per interval via a timer, or wrap main.py in a loop script.

[Unit]
Description=MyWeb Scraper (Playwright tips bot)
After=network.target

[Service]
Type=oneshot
User=gk
WorkingDirectory=/home/gk/myweb_scraper
Environment=PYTHONUNBUFFERED=1
Environment=MATCHBOT_OUTPUT_DIR=/home/gk/matchbot
ExecStart=/home/gk/myweb_scraper/venv/bin/python3 /home/gk/myweb_scraper/main.py
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
