# myweb_scraper.timer.example
# Example timer: run scraper daily at 09:30 (adjust OnCalendar as needed).
# Install:
#   sudo cp myweb_scraper.service.example /etc/systemd/system/myweb_scraper.service
#   sudo cp myweb_scraper.timer.example /etc/systemd/system/myweb_scraper.timer
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now myweb_scraper.timer

[Unit]
Description=Daily MyWeb Scraper run

[Timer]
OnCalendar=*-*-* 09:30:00
Persistent=true

[Install]
WantedBy=timers.target
